# etc/rc.local
#!/bin/bash

# GPIB board driver (tnt4882 is the NI pci card; gpib_bitbang raspi3 GPIO line hack)
#GPIB_DRIVER=tnt4882
GPIB_DRIVER=gpib_bitbang

/usr/sbin/modprobe $GPIB_DRIVER
/usr/sbin/gpib_config -t $GPIB_DRIVER --minor 0
# Without the following sleep, the owners & permissions go wrong...
sleep 2
# /usr/bin/chgrp gpib /dev/gpib*
/usr/bin/chgrp gpio /dev/gpib*
# After this, accessible to users in gpib group + everyone (reduce for security)
/usr/bin/chmod 777 /dev/gpib*
exit 0
